# src/RaspberryPi/etc/rc.local
# Add this in to deal with a bug in iptables not loading properly
sudo iptables -L >/dev/null
sudo systemctl stop openvpn
sudo systemctl start openvpn
